# packages/ns-api/files/uci-defaults/19-ns-api.wizard
#!/bin/sh

#
# Copyright (C) 2025 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-2.0-only
#

# This sets the wizard as completed if the system has been already configured, this can happen if the file
# /etc/config/dpi does not exist, to do so, this script must run before the 20_dpi uci-default
if [ ! -f /etc/config/dpi ]
then
  # system is already configured
  uci set ns-wizard.config.complete=1
  uci set ns-wizard.config.password_changed=1
  uci commit ns-wizard
fi
