# packages/ns-api/files/uci-defaults/99-ns-api.wireguard
#!/bin/sh

#
# Copyright (C) 2025 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-2.0-only
#

set -e

/usr/libexec/wireguard-migrate
