# files/etc/uci-defaults/90-nextsec-root-password
# Change root password, if needed
[ "$(grep '^root' /etc/shadow)" = "root:::0:99999:7:::" ] || exit 0

echo 'root:$6$EtySwJxJPxJiWXnw$RfgZOdGa7gBQrLagSGUD3U3bm/WeYpf0ggWpXZav.S9/2zd718maC2PKlaRTYzmUBplvbfy0Tbr4EeE5gFvyA0:19319:0:99999:7:::' > /etc/shadow.tmp

# Preserve everything else
grep -v '^root:' /etc/shadow >> /etc/shadow.tmp
mv /etc/shadow.tmp /etc/shadow
