# files/etc/uci-defaults/95-nethsec-migrate
archive="/usr/share/migration/export.tar.gz"
if [ -f "${archive}" ]; then
    /usr/sbin/ns-import $archive | tee /var/log/migration.log

    if [ $? -eq 0 ]; then
        mv "${archive}" "${archive}.done"
    fi
fi
