# files/etc/uci-defaults/99-nextsec-netifyd
if [ -f /usr/lib/libnetify-flow-actions.so.0.0.0 ] && ! grep -q np-nfa /etc/netifyd.conf; then
    (
        echo
        echo "[plugin_detections]"
        echo "np-nfa = /usr/lib/libnetify-flow-actions.so.0.0.0"
    ) >> /etc/netifyd.conf
fi
