# files/etc/uci-defaults/99-nethsec-remove-scheduled-upgrade
/usr/libexec/ns-api/schedule-system-update remove
