# files/etc/uci-defaults/99-nextsec-nginx
uci -q batch << EOI
set nginx._lan.error_log=syslog:server=unix:/dev/log
set nginx._lan.access_log=syslog:server=unix:/dev/log
commit system
EOI
