# files/etc/uci-defaults/99-nextsec-rsyslog
/etc/init.d/rsyslog stop
/etc/init.d/rsyslog disable
