# files/etc/uci-defaults/99-nethsec-backupurl
[ "$(uci get ns-plug.config.backup_url)" != "https://backupd.nethesis.it/backup" ] && exit 0

uci -q batch << EOI
set ns-plug.config.backup_url="https://backupd.nethesis.it"
commit rpcd
EOI
