# files/etc/uci-defaults/99-nethsec-restore-extra-packages
/etc/init.d/ns-restore-extra-packages start
