# files/etc/uci-defaults/99-nextsec-dnsdist
conf="/etc/dnsdist.conf"
if [ -s "$conf" ]; then
    exit 0
fi

mkdir -p "$conf.d"
echo "controlSocket('127.0.0.1:5199')" > $conf
echo "makeKey()" | dnsdist  -l 127.0.0.1:5199  | grep setKey >> $conf
echo "includeDirectory('$conf.d')" >> $conf
