# files/etc/uci-defaults/99-nethsec-pppoe
/usr/share/ns-api/remove-pppoe-keepalive
