# files/etc/uci-defaults/99-nethsec-nginx
uci -q batch << EOI
set nginx._lan.error_log=syslog:server=unix:/dev/log,nohostname
set nginx._lan.access_log=syslog:server=unix:/dev/log,nohostname
commit nginx
EOI

# make sure nginx does not listen on port 80 to avoid conflicts with acme.sh
uci -q delete nginx._redirect2ssl

exit 0
