# files/etc/uci-defaults/99-nextsec-ts-dns
/etc/init.d/adblock disable
/etc/init.d/ts-dns enable
