# files/etc/uci-defaults/99-nethsec-root-bash
sed -i "/^root:/s/\/ash$/\/bash/" /etc/passwd
