# files/etc/uci-defaults/99-nextsec-sshd
service sshd stop || return 0
service sshd disable
