# files/etc/uci-defaults/99-nextsec-sshd
service sshd stop || return 0
service sshd disable
sed -i "s/#Port 22/Port 0/g" /etc/ssh/sshd_config
