# files/etc/uci-defaults/99-nethsec-webui
uci batch <<EOI
del_list nginx._lan.include=restrict_locally
commit
EOI
