# packages/ns-objects/files/addrwatcher.init
#!/bin/sh /etc/rc.common

START=90

USE_PROCD=1
NAME=addrwatcher
PROG=/usr/sbin/addrwatcher

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_close_instance
}

reload_service() {
	procd_send_signal $PROG
}
