# packages/ns-objects/files/mac2ip.init
#!/bin/sh /etc/rc.common

#
# Copyright (C) 2022 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-2.0-only
#

START=90

USE_PROCD=1
NAME=mac2ip
PROG=/usr/sbin/mac2ip
CONFIG=/var/etc/mac2ip.json

start_service() {
	# init config
	if [ ! -f "$CONFIG" ]; then
		echo null > "$CONFIG"

	fi
	procd_open_instance
	procd_set_param command "$PROG"
	procd_close_instance
}

reload_service() {
	procd_send_signal $PROG SIGUSR1
}
