# packages/ns-threat_shield/files/threat_shield.init
#!/bin/sh /etc/rc.common

START=20

start() {
    /usr/sbin/ts
}
